# thirdparty/videoserver/zimbra-videoserver/debian/zimbra-videoserver.videoserver.service
[Unit]
Description=Zextras Videoserver
Wants=network.target

[Service]
Type=simple
ExecStart=/opt/zimbra/common/bin/janus
User=videoserver
Group=videoserver
Restart=on-failure
LimitNOFILE=65536
TasksMax=infinity

[Install]
WantedBy=multi-user.target
